# docker-zfs-plugin.socket
[Socket]
ListenStream=/run/docker/plugins/zfs.sock

[Install]
WantedBy = sockets.target
